# pfuschplay.service
[Unit]
Description=UART TFT Implementation for Klipper
After=multi-user.target
 
[Service]
ExecStart=python3 /home/pi/pfuschplay/main.py
 
[Install]
WantedBy=multi-user.target